# crontab
# Optional: Cron schedule for automated scraping
# Uncomment and customize as needed

# Run scraper daily at 2 AM
0 2 * * * cd /app && python -m scraper.main >> /var/log/scraper.log 2>&1

# Extract data after scraping (2:30 AM)
30 2 * * * cd /app && python -m scraper.extract_data >> /var/log/extractor.log 2>&1

# Generate reports (3 AM)
0 3 * * * cd /app && python -m scraper.report >> /var/log/reports.log 2>&1

# Run orchard analysis weekly (Mondays at 4 AM)
0 4 * * 1 cd /app && python -m scraper.orchard_analysis >> /var/log/orchard.log 2>&1

# Generate Word report weekly (Mondays at 4:15 AM)
15 4 * * 1 cd /app && python -m scraper.generate_agent_report >> /var/log/agent_report.log 2>&1

# Keep an empty line at the end for cron
